# tools/com.genesis.splicesorter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.genesis.splicesorter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/brawny/genesis/tools/splice_sorter.py</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/brawny/Splice/sounds/packs</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>3</integer>

    <key>StandardOutPath</key>
    <string>/Users/brawny/genesis/tools/splice_sort_launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/brawny/genesis/tools/splice_sort_launchd.log</string>
</dict>
</plist>
